# etc/service/frpc/run
#!/bin/sh -e

sv start dockergen || exit 1

exec /usr/local/bin/frpc -c /etc/frpc.toml 2>&1

